# supervisord.conf
[supervisord]
nodaemon=true

[program:docker-gen]
priority=1
command=/etc/docker-gen -config /etc/docker-gen.cfg

[program:haproxy]
priority=999
command=haproxy -f /etc/haproxy/haproxy.cfg